# bin/raceai.service
[Unit]
Description=Init Raceai Environment
After=multi-user.target
[Service]
Type=simple
ExecStart=/raceai/raceai.rc
WorkingDirectory=/raceai
Restart=always
User=pi
Group=pi
[Install]
WantedBy=multi-user.target
